# bridge/deploy/aptracker-bridge.service
[Unit]
Description=ApTracker Archipelago bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust User and paths to match your Raspberry Pi setup.
User=pi
WorkingDirectory=/home/pi/ApTracker/bridge
EnvironmentFile=/home/pi/ApTracker/bridge/.env
ExecStart=/home/pi/ApTracker/bridge/.venv/bin/python -m aptracker_bridge
Restart=always
RestartSec=5
# Keep logs in the journal (view with: journalctl -u aptracker-bridge -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
